# rpmbuild-digioceanfs/BUILD/digioceanfs_gui-3.8.0/install-conf/digioceanfs-gui.service
[Unit]
Description=digioceanfs-gui Service
After=graphical.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/digioceanfs-gui start
Restart=on-abort


[Install]
WantedBy=multi-user.target
