# rpmbuild-digioceanfs/BUILD/digioceanfs_manager-3.8.0/conf/etc/cron.d/logrotate
# Run system wide raid-check once a week on Sunday at 1am by default
*/10 * * * * root /etc/cron.daily/logrotate

