# rpmbuild-digioceanfs/BUILD/digioceanfs_manager-3.8.0/core/scripts/60-blockdev-hotplug.rules
KERNEL=="sd[a-z]*",ACTION=="add",RUN+="/usr/bin/python /usr/bin/disk_plug add /digioceanfs/wwn-$env{ID_WWN_WITH_EXTENSION}"
KERNEL=="sd[a-z]*",ACTION=="remove",RUN+="/usr/bin/python /usr/bin/disk_plug remove /digioceanfs/wwn-$env{ID_WWN_WITH_EXTENSION}"