# rpmbuild-digioceanfs/BUILD/digioceanfs_manager-3.8.0/conf/etc/udev/rules.d/10-my.rules
ACTION=="add",KERNEL=="sd*",RUN+="/usr/local/digioceanfs_manager/monitor_reporter/udev_monitor.py %k up"
ACTION=="remove",KERNEL=="sd*",RUN+="/usr/local/digioceanfs_manager/monitor_reporter/udev_monitor.py %k down"
