# rpmbuild-digioceanfs/BUILD/digioceanfs_client-3.8.0/digioceanfs-client.service
[Unit]
Description=digioceanfs-client Service
After=network.target rpcbind.service
Before=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/digioceanfs-client start
Restart=on-abort


[Install]
WantedBy=multi-user.target