# rpmbuild-digioceanfs/BUILD/digioceanfs_manager-3.8.0/conf/etc/systemd/system/multi-user.target.wants/node-manager.service
[Unit]
Description=node-manager Service
After=network.target rpcbind.service
Before=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/node-manager start
Restart=on-abort


[Install]
WantedBy=multi-user.target
